# sources/meta-ea/recipes-ea/ea-files/ea-files/other/gnssdata.service
[Unit]
Description= Colect GNSS data
After=basic.target

[Service]
Type=simple
ExecStart=/home/root/gnssdata/gnssdata &
Restart=always

[Install]
WantedBy=default.target

